# roles/goldentooth.setup_slurm_exporter/templates/slurm-exporter.service.j2
[Unit]
Description=Prometheus Slurm Exporter
Documentation=https://github.com/vpenso/prometheus-slurm-exporter
After=network.target slurmd.service slurmctld.service
Wants=network.target

[Service]
Type=simple
User=slurm-exporter
Group=slurm-exporter
ExecStart=/usr/local/bin/prometheus-slurm-exporter \
  -web.listen-address={{ ansible_default_ipv4.address }}:{{ prometheus_slurm_exporter.port }} \
  -web.log-level=info

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log/slurm-exporter

# Resource limits
LimitNOFILE=8192
MemoryMax=128M

# Restart settings
Restart=always
RestartSec=5
TimeoutStopSec=20

# Environment
Environment=PATH=/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/var/lib/slurm-exporter

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slurm-exporter

[Install]
WantedBy=multi-user.target